# scheduling/launchd/com.USER.AGENT.session.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd job: fires one scheduled autonomous session of your loopling.
  Substitute {{USER}}, {{AGENT_NAME}}, the Label suffix, the time, and the prompt.
  Install:
    cp this  ~/Library/LaunchAgents/com.{{USER}}.{{AGENT_NAME}}.sessiona.plist
    launchctl load ~/Library/LaunchAgents/com.{{USER}}.{{AGENT_NAME}}.sessiona.plist
  (see ../SETUP.md for enable/disable/permanent-disable commands)
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.{{USER}}.{{AGENT_NAME}}.sessiona</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/{{USER}}/claude-apps/{{AGENT_NAME}}/cron-wrapper.sh</string>
        <string>Session A</string>
        <string>You are {{AGENT_DISPLAY}}, running as a SCHEDULED job (no live Telegram tool in this mode). Read your soul first: ~/claude-apps/{{AGENT_NAME}}/.claude/CLAUDE.md. If you use the browser, pass `--session __CRON_SESSION__` (a unique id is substituted in) so you never collide with the live session. To message the owner, run: python3 voice/tg_send.py "your update" (one-way push — never opens the live polling connection). Then run Session A: [DESCRIBE THIS SESSION'S JOB, e.g. "morning pass — search the wiki, gather what's relevant to my purpose, file the findings, and send a short Telegram briefing via tg_send.py."]</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/{{USER}}</string>
        <key>PATH</key>
        <string>/Users/{{USER}}/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/{{USER}}/.browser-use-env/bin</string>
        <!-- So the scheduled agent can report via voice/tg_send.py (one-way push). -->
        <key>LOOPLING_BOT_TOKEN</key>
        <string>{{BOT_TOKEN}}</string>
        <key>LOOPLING_OWNER_CHAT_ID</key>
        <string>{{OWNER_CHAT_ID}}</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/{{USER}}/claude-apps/{{AGENT_NAME}}/data/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/{{USER}}/claude-apps/{{AGENT_NAME}}/data/launchd.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
